# api/supervisord.conf
[supervisord]
nodaemon=true

[program:sensu-api]
command=/etc/init.d/sensu-api start
